# question-service.conf.j2
[program:question-service]
command=/bin/sh run.sh
directory={{ path }}
stopasgroup=true
